# gateone/templates/gateone.ini
[program:gateone]
command=gateone
    --settings_dir=/opt/gateone/conf.d
    --log_file_prefix={{ gateone_path }}/logs/gateone.log